# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Description=
Exec=conky -c "/home/solus/.conky/Titus.conkyrc" &
