# .config/autostart/devil.desktop
[Desktop Entry]
Type=Application
Name=devil
Description=
Exec=devilspie2
