# .config/autostart/ibusbamboo.desktop
[Desktop Entry]
Type=Application
Name=ibusbamboo
Description=
Exec=ibus-daemon --xim &
