# .config/autostart/Set keypress delay and rate.desktop
[Desktop Entry]
Type=Application
Name=Set keypress delay and rate
Description=
Exec=xset r rate 200 70
